# udp-forwarder.service
[Unit]
Description=UDP Forwarder Service
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/udp-forwarder -config /etc/udp-forwarder/config.yaml
Restart=always
RestartSec=5s

# Allow binding to privileged ports (<1024) and raw sockets (for transparent forwarding) even when running as nobody
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
NoNewPrivileges=true

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
